# gunicorn/supervisord.conf
[program:t4j]
command='gunicorn -c gunicorn_config.py t4j.wsgi'
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/t4j.err.log
stdour_logfile=/var/log/gunicorn/t4j.out.log
